# crontab
0 0 * * * /usr/local/bin/python /app/main.py >> /app/cron.log 2>&1
